# deploy/sannsvar.service
[Unit]
Description=Sannsvar Stremio Addon
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sannsvar
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/sannsvar/.env

[Install]
WantedBy=multi-user.target
